# ansible/roles/logclean-job/templates/logclean.timer.j2
[Unit]
After=multi-user.target

[Timer]
OnCalendar={{ logclean_timer_frequency }}
Persistent=yes
